# templates/default/upstart-vpnc.conf.erb
start on (started networking)

stop on starting shutdown

console output

exec /usr/sbin/vpnc<%= '-connect' if node['platform_family'] == 'debian' %> --no-detach --pid-file /var/run/vpnc/pid --ifmode tun --non-inter --local-port 0 --dpd-idle 0

# Restart the process if it dies with a signal
# or exit code not given by the 'normal exit' stanza.
respawn
