# templates/default/upstart-vpnc.conf.erb
# /etc/init/vpnc.conf - VPN concentrator for Cisco IPSec

description "VPN concentrator for Cisco IPSec"
start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit unlimited

script
  exec /usr/sbin/vpnc<%= '-connect' if node['platform_family'] == 'debian' %> --no-detach --pid-file /var/run/vpnc/pid --non-inter --dpd-idle 0
end script
